# aria2/crontab
0 0 * * 0 /bin/ash /aria2/autoUpdateTracker.sh
